# deploy/email-engine.service
[Unit]
Description=Email Engine - PowerMTA + MailWizz Infrastructure Manager
After=network.target pmta.service

[Service]
Type=simple
User=email-engine
Group=email-engine
WorkingDirectory=/opt/email-engine
Environment="PATH=/opt/email-engine/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/email-engine/.env
ExecStart=/opt/email-engine/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email-engine

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/email-engine /etc/pmta
ProtectHome=true

[Install]
WantedBy=multi-user.target
